# deploy/piloop-bridge.service
[Unit]
Description=pi-loop BLE console bridge (PTY <-> Nordic UART)
# BlueZ for GATT/advertising. Wraps the agent loop so its console streams over
# BLE to tab-device. NOT a boot dependency of anything — BlueZ starts late.
After=bluetooth.target piloop-agent.service
Wants=bluetooth.target

[Service]
# This service IS the running loop. It forks loop.py as its child under the PTY
# and is the unit the host's status/logs/deploy commands target.
EnvironmentFile=-/etc/piloop/agent.env
ExecStart=/usr/bin/python3 /opt/piloop/bridge/bridge.py -- /usr/bin/python3 /opt/piloop/agent/loop.py
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
